# mysqld_exporter/deb/control/postinst
#!/bin/bash

is_systemd=0

if [ -d /run/systemd/system ]; then
   is_systemd=1
fi

if [ $is_systemd -eq 1 ]
then
	systemctl daemon-reload && systemctl enable mysqld-exporter
fi
