# mysqld_exporter/deb/control/prerm
#!/bin/bash

invoke-rc.d mysqld-exporter stop
update-rc.d -f mysqld-exporter remove
